# secret-scanner.service
[Unit]
Description=Confluence Secret Scanner Service
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/opt/confluence-automation
Environment="PATH=/opt/confluence-automation/.venv/bin"
ExecStart=/opt/confluence-automation/.venv/bin/python /opt/confluence-automation/secret_scanner.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
